# docker/supervisord.conf
; =============================================================================
; supervisord.conf — runs php-fpm + queue worker + scheduler in one container
; Use this as an alternative to the separate queue/scheduler services in
; docker-compose.yml (mount it as /etc/supervisord.conf in the app service).
; =============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; PHP-FPM
; =============================================================================
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/www/html/storage/logs/php-fpm.log
stderr_logfile=/var/www/html/storage/logs/php-fpm-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; =============================================================================
; Queue worker — processes jobs from all named queues
; =============================================================================
[program:queue-worker]
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --queue=default,high,low
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
startretries=5
numprocs=2
process_name=%(program_name)s_%(process_num)02d
priority=10
stdout_logfile=/var/www/html/storage/logs/queue-worker.log
stderr_logfile=/var/www/html/storage/logs/queue-worker-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; =============================================================================
; Scheduler — runs `php artisan schedule:run` every 60 seconds
; Alpine images have no cron, so we use a shell loop.
; =============================================================================
[program:scheduler]
command=sh -c "while true; do php /var/www/html/artisan schedule:run --verbose --no-interaction >> /var/www/html/storage/logs/scheduler.log 2>&1; sleep 60; done"
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
startretries=3
priority=20
stdout_logfile=/var/www/html/storage/logs/scheduler-supervisor.log
stderr_logfile=/var/www/html/storage/logs/scheduler-supervisor-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
